# supervisor_conf/udp_logger-web_server.conf
[program:udp_logger-web_server]
command=python3 /root/udp_logger/web_server.py
autorestart=true