# supervisor_conf/udp_logger-udp_server.conf
[program:udp_logger-udp_server]
command=python3 /root/udp_logger/udp_server.py
autorestart=true